# systemd/mb8600-watchdog.service
[Unit]
Description=MB8600 Watchdog - Watchdog Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/opt/mb8600-watchdog/bin/watchdog --config /etc/mb8600-watchdog/config.json
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
